# apod.service
[Unit]
Description=apod wallpapers
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/adi-bg

[Install]
WantedBy=default.target
