# apod.service
[Unit]
Description=apod wallpapers
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/sid-bg

[Install]
WantedBy=default.target
