# systemd/sid-bg.timer
[Unit]
Decription=Run sid-bg daily at midnight or on boot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target