# systemd/sid-bg.service
[Unit]
Description=Download newest apod image
After=multi-user.target

[Service]
Type=simple
ExecStart=sid-bg

[Install]
WantedBy=default.target